# systemd/tradingbot@.service.d/onfailure.conf
[Unit]
OnFailure=tradingbot-alert@%i.service
