# systemd/tradingbot-web.service
[Unit]
Description=Trading Bot Web Dashboard
After=network.target
PartOf=tradingbot.target

[Service]
Type=simple
User=DEIN_USER
WorkingDirectory=DEIN_BOTDIR
EnvironmentFile=DEIN_BOTDIR/.env
ExecStart=DEIN_BOTDIR/botvenv/bin/python DEIN_BOTDIR/web/app.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=tradingbot.target
