# systemd/tradingbot.target
[Unit]
Description=Alle Trading Bot Instanzen
After=network.target

[Install]
WantedBy=multi-user.target
